# deploy/ansible/templates/ivy-vine.service.j2
[Unit]
Description=Ivy Vine Daemon
After=network-online.target docker.service postgresql.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=ivy
Group=ivy
ExecStart=/usr/local/bin/ivy-vine --config /etc/ivy/config.yaml
Restart=always
RestartSec=10
LimitNOFILE=65536

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ivy-vine

[Install]
WantedBy=multi-user.target
